# pubnub-temp-site/deploy/publish.service
[Unit]
Description=PubNub test publisher
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pubnub
# Use EnvironmentFile or pass env vars via systemd drop-in for secrets
Environment=PUB_KEY=pub-xxx
Environment=SUB_KEY=sub-xxx
ExecStart=/usr/bin/node /opt/pubnub/publish_test.js my-channel --loop --interval 5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
